# autostart_services/fixf.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/fixf.sh

[Install]
WantedBy=multi-user.target
